# scripts/company-portal.service
[Unit]
Description=Company Portal API Server
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/company-portal/backend
ExecStartPre=/bin/bash -c 'fuser -k 3001/tcp 2>/dev/null || true'
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=company-portal
Environment=NODE_ENV=production
# Give the process time to start
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
